# scratchy-app.service
[Unit]
Description=Scratchy ASP.NET Core Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/scratchy
ExecStart=/var/www/scratchy/Scratchy
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=scratchy-app
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

[Install]
WantedBy=multi-user.target
